# system/systemd/user/mugen-event-notifier.service
[Unit]
Description=mugen-shell calendar event notifications

[Service]
Type=oneshot
ExecStart=%h/.config/quickshell/mugen-shell/scripts/notify-events.py
