# system/systemd/user/yura-voice.service
[Unit]
Description=Yura voice input daemon (wake word -> STT -> chat -> TTS)
# Needs the compositor for qs IPC and PipeWire for the mic.
After=graphical-session.target voicevox-engine.service mugen-ai.service
Wants=voicevox-engine.service
PartOf=graphical-session.target

[Service]
WorkingDirectory=%h/mugen-shell/voice
# Custom model trained in voice/train; scores saturate near 1.0 on true
# positives, so 0.7 trims borderline false wakes at ~no recall cost.
Environment=YURA_WAKEWORD=%h/mugen-shell/voice/models/hey_yura.onnx
Environment=YURA_WAKE_THRESHOLD=0.7
ExecStart=%h/mugen-shell/voice/.venv/bin/python %h/mugen-shell/voice/yurad.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
